# base/etc/service/cron/run
#!/bin/sh -e

exec /usr/sbin/cron -f 2>&1

