# x11/etc/service/xpra-proxy-html5/run
#!/bin/sh -e

# include configuration file that should do all checks and setup variables
. ./configure.sh

echo '######################################################################################################' #| tee $LOG_DIR/xpra-proxy.log
echo '######################################### STARTED XRPA Proxy #########################################' #| tee $LOG_DIR/xpra-proxy.log
# Force "threaded" (in-process) proxy instances instead of the default per-client child
# process. When the HTML5 client asks the proxy to START a new session (action=start), xpra
# 5.x hands the client socket over to a freshly-spawned child proxy process via fd-passing.
# That handover only works for ssl/wss/ssh client sockets; for a plain "ws" socket - which
# is what we get because FRP terminates TLS upstream and forwards plain HTTP - the handover
# fails ("network IO threads have failed to terminate") and the browser hangs forever at
# "Sending handshake" / "WebSocket connection established". Threaded instances relay the
# connection in-process with no handover, so starting sessions on demand works behind the
# reverse proxy. (This regressed in xpra 5.x; 4.4.x relayed differently.)
export XPRA_PROXY_INSTANCE_THREADED=1

# tcp-auth=allow issues a password challenge and accepts any answer (the container is
# single-user; real access control is the front proxy/FRP + the network). For this to work
# with a modern HTML5 client (v14+), the server must NOT challenge with the insecure "xor"
# digest: since v14 the client refuses xor over a connection it considers unencrypted (and
# behind a TLS-terminating reverse proxy it cannot tell the hop is encrypted), dropping the
# WebSocket. xpra's choose_digest() falls back to xor whenever the client offers none of
# hmac+sha512/384/256/224 - skipping the hmac/hmac+md5 every client also offers - so the
# Dockerfile patches choose_digest to prefer those safe hmac digests before xor. The server
# then challenges with hmac (never xor) and the modern client accepts it.
exec chpst -u $USER_NAME xpra proxy :1000 \
						--tcp-auth=allow \
						--bind-tcp=0.0.0.0:${XPRA_PROXY_HTML5_PORT} \
						--html=/usr/share/xpra/www \
						--clipboard=GDK --clipboard-direction=both \
						--bell=off --daemon=no >> $LOG_DIR/$HOSTNAME-xpra-proxy.log 2>&1


