# x11/etc/service/dbus/run
#!/bin/sh -e

mkdir -p /var/run/dbus

# clear any old dbug pids
rm -f /var/run/dbus/pid 2> /dev/null

exec /usr/bin/dbus-daemon --system --nofork 2>&1

