# x11/etc/service/dbus/run
#!/bin/sh -e

mkdir -p /run/dbus /var/run/dbus

# Container restart preserves the writable runtime filesystem, so clear any
# stale dbus state left behind by the previous boot before starting a new bus.
rm -f /run/dbus/pid /var/run/dbus/pid /run/dbus/system_bus_socket /var/run/dbus/system_bus_socket 2> /dev/null

exec /usr/bin/dbus-daemon --system --nofork 2>&1
