# x11/etc/service/dbus/run
#!/bin/sh -e

mkdir -p /var/run/dbus

exec /usr/bin/dbus-daemon --system --syslog-only --nofork 2>&1

