# base/etc/service/rsyslog/run
#!/bin/sh -e

exec /usr/sbin/rsyslogd -n 2>&1

